# ops/supervisord/helloworld.conf
[program:helloworld]
command = /vagrant/ops/bash/start_gunicorn.sh